# accuknox-files/v0.2.13/usr/lib/systemd/system/accuknox-shared-informer-agent.service
[Unit]
Description=Shared-Informer

[Service]
User=root
KillMode=process
WorkingDirectory=/opt/accuknox-shared-informer-agent/
ExecStart=/opt/accuknox-shared-informer-agent/sia

[Install]
WantedBy=multi-user.target
